# etc/supervisor/conf.d/im.conf
[program:skype]
command=/usr/bin/skype
autostart=true
autorestart=true

[program:pidgin]
command=/usr/bin/pidgin
autostart=true
autorestart=true

[group:im]
programs=pidgin,skype
